# rc-service
#!/sbin/openrc-run

LYNXIE_USER="uncertainty"
LYNXIE_CMD="/home/uncertainty/lynxie"
APP_NAME="lynxie"

command="$LYNXIE_CMD"
command_args="-prefix=> -datapath=_data -token=Nope"
command_user="$LYNXIE_USER"
pidfile="/var/run/${APP_NAME}.pid"
output_log="/var/log/${APP_NAME}.log"
error_log="/var/log/${APP_NAME}_error.log"

depend() {
        need net
}

start() {
        ebegin "Starting ${APP_NAME}"

        # The bot also searches for the bot token in the env
        # export TOKEN=""

        start-stop-daemon --start               \
                          --background          \
                          --user "$LYNXIE_USER" \
                          --make-pidfile        \
                          --pidfile "$pidfile"  \
                          --exec "$command" -- $command_args >> "$output_log" 2>> "$error_log"

        eend $?
}

stop() {
        ebegin "Stopping ${APP_NAME}"

        start-stop-daemon --stop --pidfile "$pidfile"

        eend $?
}
